# debian/funkeymonkey-pyrainput.postinst
#!/bin/sh

systemctl enable pyrainput
systemctl start pyrainput
