# debian/funkeymonkey-pyrainput.prerm
#!/bin/sh

systemctl stop pyrainput
systemctl disable pyrainput
